# virtualhere_usb_server/rootfs/etc/services.d/virtualhere/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

CONFIG_DIR="/addon_configs/${HOSTNAME}"
mkdir -p "$CONFIG_DIR"
CONFIG_FILE="$CONFIG_DIR/config.ini"
BINARY_DIR="/usr/local/bin"

VH_PID=""
SHUTTING_DOWN="false"

map_arch() {
    local machine
    machine="$(uname -m)"

    case "$machine" in
        x86_64) echo "vhusbdx86_64" ;;
        i386|i486|i586|i686) echo "vhusbdi386" ;;
        aarch64|arm64) echo "vhusbdarm64" ;;
        armv7l|armv7*) echo "vhusbdarm" ;;
        armv6l|armhf|arm*) echo "vhusbdarm" ;;
        *)
            bashio::log.fatal "Unsupported architecture: $machine"
            exit 1
            ;;
    esac
}

write_config() {
    local raw
    raw="$(bashio::config 'config_ini' || true)"

    if [ -n "$raw" ]; then
        printf '%s\n' "$raw" > "$CONFIG_FILE"
        bashio::log.info "Wrote config.ini from app options to $CONFIG_FILE"
    elif [ ! -f "$CONFIG_FILE" ]; then
        : > "$CONFIG_FILE"
        bashio::log.info "Created empty config.ini at $CONFIG_FILE"
    else
        bashio::log.info "Using existing config.ini at $CONFIG_FILE"
    fi
}

download_binary() {
    local filename url target
    filename="$(map_arch)"
    url="https://www.virtualhere.com/sites/default/files/usbserver/${filename}"
    target="${BINARY_DIR}/${filename}"

    if [ ! -x "$target" ]; then
        bashio::log.info "Downloading VirtualHere server binary: ${filename}"
        curl -fsSL "$url" -o "$target"
        chmod 0755 "$target"
    else
        bashio::log.info "Using cached VirtualHere server binary: ${filename}"
    fi

    printf '%s\n' "$target"
}

usb_device_matches() {
    local path="$1"
    local vendor="$2"
    local product="$3"
    local sysfs_path="/sys/bus/usb/devices/${path}"
    local actual_vendor=""
    local actual_product=""

    if [ ! -d "$sysfs_path" ]; then
        bashio::log.warning "USB device path not found: ${sysfs_path}"
        return 1
    fi

    if [ -n "$vendor" ]; then
        if [ ! -f "${sysfs_path}/idVendor" ]; then
            bashio::log.warning "USB device ${path} has no idVendor file"
            return 1
        fi

        actual_vendor="$(cat "${sysfs_path}/idVendor")"

        if [ "$actual_vendor" != "$vendor" ]; then
            bashio::log.warning "USB device ${path} vendor mismatch: expected ${vendor}, got ${actual_vendor}"
            return 1
        fi
    fi

    if [ -n "$product" ]; then
        if [ ! -f "${sysfs_path}/idProduct" ]; then
            bashio::log.warning "USB device ${path} has no idProduct file"
            return 1
        fi

        actual_product="$(cat "${sysfs_path}/idProduct")"

        if [ "$actual_product" != "$product" ]; then
            bashio::log.warning "USB device ${path} product mismatch: expected ${product}, got ${actual_product}"
            return 1
        fi
    fi

    return 0
}

reset_usb_device() {
    local path="$1"
    local vendor="$2"
    local product="$3"

    if [ -z "$path" ]; then
        bashio::log.warning "Skipping USB reset entry with empty path"
        return 0
    fi

    bashio::log.info "Preparing USB reset for ${path}"

    if ! usb_device_matches "$path" "$vendor" "$product"; then
        bashio::log.warning "Skipping USB reset for ${path}"
        return 0
    fi

    bashio::log.info "Resetting USB device ${path}"

    if ! printf '%s\n' "$path" > /sys/bus/usb/drivers/usb/unbind; then
        bashio::log.warning "Failed to unbind USB device ${path}"
        return 0
    fi

    sleep 5

    if ! printf '%s\n' "$path" > /sys/bus/usb/drivers/usb/bind; then
        bashio::log.warning "Failed to bind USB device ${path}"
        return 0
    fi

    bashio::log.info "USB reset completed for ${path}"
}

reset_usb_devices() {
    local phase="$1"
    local count
    local index
    local path
    local vendor
    local product
    local reset_on_start
    local reset_on_stop

    count="$(bashio::config 'usb_reset_devices | length' 2>/dev/null || echo 0)"

    if [ "$count" = "null" ] || [ "$count" -eq 0 ]; then
        bashio::log.info "No USB reset devices configured for ${phase}"
        return 0
    fi

    bashio::log.info "Checking USB reset devices for ${phase}"

    for index in $(seq 0 $((count - 1))); do
        path="$(bashio::config "usb_reset_devices[${index}].path" || true)"
        vendor="$(bashio::config "usb_reset_devices[${index}].vendor" || true)"
        product="$(bashio::config "usb_reset_devices[${index}].product" || true)"
        reset_on_start="$(bashio::config "usb_reset_devices[${index}].reset_on_start" || echo false)"
        reset_on_stop="$(bashio::config "usb_reset_devices[${index}].reset_on_stop" || echo false)"

        case "$phase" in
            start)
                if [ "$reset_on_start" = "true" ]; then
                    reset_usb_device "$path" "$vendor" "$product"
                else
                    bashio::log.info "USB reset on start disabled for ${path}"
                fi
                ;;
            stop)
                if [ "$reset_on_stop" = "true" ]; then
                    reset_usb_device "$path" "$vendor" "$product"
                else
                    bashio::log.info "USB reset on stop disabled for ${path}"
                fi
                ;;
        esac
    done
}

shutdown() {
    if [ "$SHUTTING_DOWN" = "true" ]; then
        return 0
    fi

    SHUTTING_DOWN="true"

    bashio::log.info "Stop signal received for VirtualHere USB Server"

    if [ -n "$VH_PID" ] && kill -0 "$VH_PID" 2>/dev/null; then
        bashio::log.info "Sending TERM to VirtualHere USB Server"
        kill -TERM "$VH_PID" 2>/dev/null || true
    fi

    bashio::log.info "Running USB reset after VirtualHere stop signal"
    reset_usb_devices stop

    if [ -n "$VH_PID" ] && kill -0 "$VH_PID" 2>/dev/null; then
        local waited=0

        while kill -0 "$VH_PID" 2>/dev/null && [ "$waited" -lt 10 ]; do
            sleep 1
            waited=$((waited + 1))
        done

        if kill -0 "$VH_PID" 2>/dev/null; then
            bashio::log.warning "VirtualHere did not stop within 10 seconds, sending KILL"
            kill -KILL "$VH_PID" 2>/dev/null || true
        fi

        wait "$VH_PID" 2>/dev/null || true
    fi

    bashio::log.info "VirtualHere USB Server stopped"
    exit 0
}

main() {
    local bin
    local exit_code=0

    trap shutdown SIGTERM SIGINT SIGHUP

    write_config
    bin="$(download_binary)"

    reset_usb_devices start

    bashio::log.info "Starting VirtualHere USB Server"

    "$bin" -c "$CONFIG_FILE" &
    VH_PID="$!"

    wait "$VH_PID" || exit_code="$?"

    if [ "$SHUTTING_DOWN" != "true" ]; then
        bashio::log.info "VirtualHere USB Server exited with code ${exit_code}"
        bashio::log.info "Running USB reset after VirtualHere process exit"
        reset_usb_devices stop
    fi

    exit "$exit_code"
}

main "$@"